# docker/nginx_supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=nginx -g "daemon off;"
autostart=true

[program:cron]
command=cron -f -L 15
autostart=true
autorestart=true
